# conf/uwsgi.conf
start on [2345]
stop on [06]

script
    cd /home/council/dcstat
    exec uwsgi --ini /etc/dcstat-uwsgi/dcstat-uwsgi.ini
end script